# home/private_Library/LaunchAgents/com.tomdale.pnpm-npm-token-cache.plist.tmpl
{{- $isWork := get . "isWork" | default false -}}
{{- $npmTokenOpRef := get . "npmTokenOpRef" | default "" -}}
{{- if and (eq .chezmoi.os "darwin") $isWork $npmTokenOpRef -}}
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.tomdale.pnpm-npm-token-cache</string>
	<key>ProgramArguments</key>
	<array>
		<string>{{ joinPath .chezmoi.homeDir ".local/libexec/pnpm-npm-token-cache" }}</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>{{ joinPath .chezmoi.homeDir "Library/Logs/pnpm-npm-token-cache.log" }}</string>
	<key>StandardErrorPath</key>
	<string>{{ joinPath .chezmoi.homeDir "Library/Logs/pnpm-npm-token-cache.error.log" }}</string>
</dict>
</plist>
{{- end }}
